# roles/geminabox/templates/geminabox.conf.j2
[program:geminabox]
directory = /root/.gem/mirror/
command = rackup
stderr_logfile = /var/log/geminabox_error.log
stdout_logfile = /var/log/geminabox.log
autostart=true
autorestart=true