# roles/uptime/templates/uptime_monitor.conf
[program:uptime_monitor]
directory = /opt/uptime
command = /usr/bin/env node monitor.js
stderr_logfile = /var/log/uptime/monitor_error.log
stdout_logfile = /var/log/uptime/monitor.log
autostart=true
autorestart=true
