# roles/uptime/templates/uptime.conf
[program:uptime]
directory = /opt/uptime
command = /usr/bin/env node app.js
stderr_logfile = /var/log/uptime/app_error.log
stdout_logfile = /var/log/uptime/app.log
autostart=true
autorestart=true
